# .docker/etc/s6/services/freeswitch/run
#!/bin/execlineb -WP
/sbin/su-exec freeswitch
/usr/bin/freeswitch -nc -nf -np -nosql -nonat -nonatmap
