# .docker/etc/s6/services/postgres/finish
#!/bin/execlineb -WP
/bin/s6-svscanctl -t /etc/s6/services
